# config/includes.chroot/etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Type=Application
Name=Plank Dock
Exec=plank
Hidden=false
X-GNOME-Autostart-enabled=true
Comment=macOS-style dock for iLinux
